# init/neutron/neutron-l3.conf
description "Neutron L3-Agent"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]

respawn

exec su -s /bin/sh -c "exec neutron-l3-agent --config-file /etc/neutron/neutron.conf \
                                             --config-file /etc/neutron/l3_agent.ini \
                                             --config-file /etc/neutron/fwaas_driver.ini > /dev/null 2>&1" root
